# {{cookiecutter.project_slug}}/ansible/systemd/files/app.service
[Unit]
Description=Job that runs the application daemon

[Service]
Type=simple
ExecStart=/srv/gunicorn.sh
